# one-api.service
# File path: /etc/systemd/system/fly-api.service
# sudo systemctl daemon-reload
# sudo systemctl start fly-api
# sudo systemctl enable fly-api
# sudo systemctl status fly-api
[Unit]
Description=Fly API Service
After=network.target

[Service]
User=ubuntu  # 注意修改用户名
WorkingDirectory=/path/to/fly-api  # 注意修改路径
ExecStart=/path/to/fly-api/fly-api --port 3000 --log-dir /path/to/fly-api/logs  # 注意修改路径和端口号
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
